# dist/debian/smtp/assets/supervisor.conf
[supervisord]
nodaemon=true

[program:rsyslog]
command=/usr/sbin/rsyslogd -n

[program:postfix]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
command=/opt/start.sh
